# addon-HassPyFrigate/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Runs the Nginx daemon
# ==============================================================================

bashio::log.info "Starting Apache..."
exec apache2
