# addon-HassPyFrigate/rootfs/etc/services.d/apache2/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Runs the Nginx daemon
# ==============================================================================

bashio::log.info "Starting Apache..."
exec /etc/init.d/apache2 start
